# deploy/mesh-deadrop.service
[Unit]
Description=Meshtastic Dead Drop Relay
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/mesh-deadrop
ExecStart=/opt/mesh-deadrop/deploy/start.sh
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
